# api-server/supervisord.conf
[supervisord]
loglevel=error
nodaemon=true

[supervisorctl]
serverurl=unix://supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[unix_http_server]
chown=deployments:deployments
file=supervisor.sock

[program:service]
command=deployctl service run
autostart=false
autorestart=false
redirect_stderr=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0

[program:proxy]
command=envoy -c envoy.yaml
autostart=false
autorestart=false
redirect_stderr=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
